# afplogin
### BEGIN INIT INFO
# Provides:          afplogin
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Data transfer from websockets to Arduino.
# Description:       Recieves data from websockets and transfers it to an Arduino.
### END INIT INFO

exec /usr/local/bin/afplogin
